# examples/client-service/tosh-client.timer
[Unit]
Description=Updates SSH config using tosh

[Timer]
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
