# examples/iptables/tosh-iptables.timer
[Unit]
Description=Updates SSH allowed IP address using tosh-iptables service

[Timer]
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
